# deployment/worksphere.service
[Unit]
Description=Worksphere Spring Boot Application
After=network.target syslog.target

[Service]
Type=simple

User=worksphere
Group=worksphere
WorkingDirectory=/opt/worksphere

EnvironmentFile=-/opt/worksphere/.env

ExecStart=/usr/bin/java -Dspring.profiles.active=production -jar /opt/worksphere/worksphere.jar

SuccessExitStatus=143
Restart=always
RestartSec=10

StandardOutput=append:/opt/worksphere/logs/worksphere-output.log
StandardError=append:/opt/worksphere/logs/worksphere-error.log

# ---- Security hardening (recommended) ----
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/worksphere

[Install]
WantedBy=multi-user.target
